# crontab
*/5 * * * * <path-to-this-repo>/run.sh > /dev/null
