# crontab
*/5 * * * * <path-to-this-repo>/run.py -k <api-key> -r <zone-id> -d <domain-name> -e <email> > /dev/null
